# playbooks/set-up-mirrors/3-install-systemd-services.d/refresh-mirror.timer
[Unit]
Description=Refresh the mirror

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
